# feeds/routing/smcroute/files/smcroute.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2010 BoluWrt.org
START=50

start() {
	smcroute -d
}

stop() {
	smcroute -k
}
