# feeds/routing/quagga/files/quagga.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 BoluWrt.org

START=60
start() {
	/usr/sbin/quagga.init start
}

stop() {
	/usr/sbin/quagga.init stop
}
